# DeployPython.service
[Unit]
Description=uWSGI instance to serve DeployPython
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/DeployPython
Environment="PATH=/home/ubuntu/DeployPython/venv/bin"
ExecStart=/home/ubuntu/DeployPython/venv/bin/gunicorn --workers 3 --bind unix:app.sock -m 007 app:app

[Install]
WantedBy=multi-user.target